# packaging/tpfan-daemon.service
[Unit]
Description=tpfan — ThinkPad Fan Control Daemon
After=systemd-modules-load.service

[Service]
Type=dbus
BusName=org.tpfan1
ExecStart=/usr/local/bin/tpfan-daemon
WatchdogSec=10
Restart=on-failure
NotifyAccess=main
Environment=TPFAN_LOG=info

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelModules=true
ProtectControlGroups=true
ConfigurationDirectory=tpfan
StateDirectory=tpfan
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
# Hinweise:
# - ProtectKernelTunables NICHT setzen: würde /proc/acpi/ibm/fan read-only machen.
# - /proc liegt außerhalb ProtectSystem=strict (read-only-Zone /usr,/boot,/etc),
#   daher ist /proc/acpi/ibm/fan ohne explizite ReadWritePaths beschreibbar.
# - MemoryDenyWriteExecute weggelassen, weil Python-/PyGObject-Code via ctypes/JIT
#   gelegentlich W+X-Mappings braucht.

[Install]
WantedBy=multi-user.target
